# workloads/kvmtool/source/guest/etc/inittab
::sysinit:sh -c "echo Hello, RISC-V from guest!; poweroff"
