# workloads/linux/rvv-bench/inittab
::sysinit:sh -c "nemu-trap 256; nemu-trap 257"
::once:sh -c "/workloads/run.sh; nemu-trap -1"
