# workloads/hello/inittab
::sysinit:/usr/bin/hello "RISC-V"
::once:/bin/nemu-halt
