# workloads/linux/coremark-pro/inittab
::sysinit:sh -c "nemu-trap 256; nemu-trap 257"
::once:sh -c "/coremark-pro/run.sh; nemu-trap -1"
