# workloads/coremark-pro/inittab
::sysinit:/coremark-pro/run.sh
::once:/bin/nemu-halt