# workloads/rvv-bench/inittab
::sysinit:/workloads/run.sh
::once:/bin/nemu-halt
