# workloads/linux/litmus-tests-riscv/inittab
::once:sh -c "/litmus-tests-riscv/run-own.sh; nemu-trap -1"
